# service/xyz.openbmc_project.provisioning.service.in
[Unit]
Description=Provisioning D-Bus Service

[Service]
ExecStart=/usr/bin/provisioningd -c /var/provisioning/provisioning.conf
Restart=always
Type=dbus
BusName=@DEFAULT_BUSNAME@
[Install]
WantedBy=network.target
